# gunicorn.conf
[program:gunicorn]
directory=/home/ubuntu/SWE-573-Project
command=/home/ubuntu/venv/bin/gunicorn --workers 3 --bind unix:/home/ubuntu/SWE-573-Project/app.sock co_learning_site.wsgi:application
autostart=true
autorestart=true
stderr_logfile=/var/log/gunicorn/gunicorn.err.log
stdout_logfile=/var/log/gunicorn/gunicorn.out.log


[group:guni]
programs:gunicorn